# deploy/systemd/velociraptor-asset-id-fetcher.service
[Unit]
Description=Velociraptor asset-id fetcher
Documentation=https://github.com/spockwall/velociraptor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ben
Group=ben
WorkingDirectory=/home/ben/velociraptor
ExecStartPre=/usr/bin/test -x /home/ben/velociraptor/target/release/asset_id_fetcher
ExecStart=/home/ben/velociraptor/target/release/asset_id_fetcher --config configs/fetcher.yaml
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
