# deploy/systemd/velociraptor-polymarket-recorder.service
[Unit]
Description=Velociraptor Polymarket recorder
Documentation=https://github.com/spockwall/velociraptor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ben
Group=ben
WorkingDirectory=/home/ben/velociraptor
ExecStartPre=/usr/bin/test -x /home/ben/velociraptor/target/release/polymarket_recorder
ExecStart=/home/ben/velociraptor/target/release/polymarket_recorder --config configs/prod/polymarket.yaml
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
