# deploy/systemd/velociraptor-polymarket-recorder.service
[Unit]
Description=Velociraptor Polymarket recorder
Documentation=https://github.com/spockwall/velociraptor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=velociraptor
Group=velociraptor
WorkingDirectory=/opt/velociraptor
ExecStartPre=/usr/bin/test -x /opt/velociraptor/target/release/polymarket_recorder
ExecStart=/opt/velociraptor/target/release/polymarket_recorder --config configs/polymarket.yaml
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
