# deploy/systemd/velociraptor-orderbook-recorder.service
[Unit]
Description=Velociraptor orderbook recorder
Documentation=https://github.com/spockwall/velociraptor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ben
Group=ben
WorkingDirectory=/home/ben/velociraptor
ExecStartPre=/usr/bin/test -x /home/ben/velociraptor/target/release/orderbook_recorder
ExecStart=/home/ben/velociraptor/target/release/orderbook_recorder --config configs/server.yaml
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
